# masakari-processmonitor/etc/systemd/system/masakari-processmonitor.service
[Unit]
Description=Masakari Processmonitor

[Service]
Type=simple
User=openstack
Group=openstack
ExecStart=bash /opt/masakari/hostmonitor/masakari-processmonitor.sh

[Install]
WantedBy=multi-user.target
